# templates/default/elasticsearch.upstart.erb
# Generated by Chef - manual changes made will be lost!
description     "Elasticsearch server instance"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# If the process quits unexpectedly trigger a respawn
respawn
respawn limit 3 30

# Increase ulimit to recommended levels
limit nofile <%= node['elasticsearch']['limits']['nofile'] %> <%= node['elasticsearch']['limits']['nofile'] %>
limit memlock <%= node['elasticsearch']['limits']['memlock'] %> <%= node['elasticsearch']['limits']['memlock'] %>

# Start the process
script
  export JAVA_HOME='<%= "#{node['java']['java_home']}" %>'
  export ES_HOME='<%= "#{node['elasticsearch']['path']['install']}/elasticsearch" %>'
  export ES_CLASSPATH=$ES_CLASSPATH:$ES_HOME/lib/*:$ES_HOME/lib/sigar/*
  export ES_HEAP_SIZE=<%= node['elasticsearch']['allocated_memory'] %>
  export ES_JAVA_OPTS="
    -server
    -Djava.net.preferIPv4Stack=true
    -Des.config=<%= @config_dir %>/elasticsearch.yml
    -Xms<%= node['elasticsearch']['allocated_memory'] %>
    -Xmx<%= node['elasticsearch']['allocated_memory'] %>
    -Xss<%= node['elasticsearch']['thread_stack_size'] %>
    <%= node['elasticsearch']['gc_settings'] %>
    <%= node['elasticsearch']['env_options'] %>
  "

  su -s /bin/dash -c "<%= node['elasticsearch']['path']['install'] %>/bin/elasticsearch -f" <%= node['elasticsearch']['user'] %>
end script